# .config/systemd/user/compton.service
[Unit]
Description=Compton
After=xorg.target

[Service]
ExecStart=compton
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
