# .config/systemd/user/wicd-tray.service
[Unit]
Description=WICD-Tray
[Service]
ExecStart=wicd-gtk -t
[Install]
WantedBy=default.target
